# etc/s6-overlay/services.d/tor/run
#!/usr/bin/execlineb -P

tor -f /etc/tor/torrc
