# etc/s6-overlay/services.d/nginx/run
#!/usr/bin/execlineb -P

# privilege escalation is being performed by nginx in nginx.conf
# s6-setuidgid http
nginx -g "daemon off;"
